# templates/etc/systemd/system/systemd.service.j2
# {{ ansible_managed }}

[Unit]
Description="{{ item.value.description | default("Service for " + item.key) }}"
{% if item.value.wants is defined %}
Wants={{ item.value.wants }}
{% endif -%}
{% if item.value.requires is defined %}
Requires={{ item.value.requires }}
{% endif -%}
{% if item.value.before is defined %}
Before={{ item.value.before }}
{% endif -%}
{% if item.value.after is defined %}
After={{ item.value.after }}
{% endif %}

[Service]
{% if item.value.type is defined %}
Type={{ item.value.type }}
{% endif -%}

{% if item.value.timeoutsec is defined %}
TimeoutSec={{ item.value.timeoutsec }}
{% endif -%}
{% if item.value.restartsec is defined %}
RestartSec={{ item.value.restartsec }}
{% endif -%}
{% if item.value.restart is defined %}
Restart={{ item.value.restart }}
{% endif -%}

{% if item.value.environment is defined %}
Environment={{ item.value.environment }}
{% endif -%}
{% if item.value.environmentfile is defined %}
EnvironmentFile={{ item.value.environmentfile }}
{% endif -%}

{% if item.value.workingdirectory is defined %}
WorkingDirectory={{ item.value.workingdirectory }}
{% endif -%}
{% if item.value.user is defined %}
User={{ item.value.user }}
{% endif -%}
{% if item.value.group is defined %}
Group={{ item.value.group }}
{% endif -%}

{% if item.value.oomscoreadjust is defined %}
OOMScoreAdjust={{ item.value.oomscoreadjust }}
{% endif -%}
{% if item.value.nice is defined %}
Nice={{ item.value.nice }}
{% endif -%}

{% if item.value.execstop is defined %}
ExecStop={{ item.value.execstop }}
{% endif -%}
{% if item.value.execstoppost is defined %}
ExecStopPost={{ item.value.execstoppost }}
{% endif -%}
{% if item.value.execreload is defined %}
ExecReload={{ item.value.execreload }}
{% endif -%}

{% if item.value.execstartpre is defined %}
ExecStartPre={{ item.value.execstartpre }}
{% endif -%}
{% if item.value.execstartpost is defined %}
ExecStartPost={{ item.value.execstartpost }}
{% endif -%}
{% if item.value.execstart is defined %}
ExecStart={{ item.value.execstart }}
{% endif -%}

{% if item.value.standardinput is defined %}
StandardInput={{ item.value.standardinput }}
{% endif -%}
{% if item.value.standardoutput is defined %}
StandardOutput={{ item.value.standardoutput }}
{% endif -%}
{% if item.value.standarderror is defined %}
StandardError={{ item.value.standarderror }}
{% endif %}

[Install]
{% if item.value.wantedby is defined %}
WantedBy={{ item.value.wantedby }}
{% endif -%}
{% if item.value.requiredby is defined %}
RequiredBy={{ item.value.requiredby }}
{% endif %}
